# logsync.timer
[Unit]
Description="Inicia o logsync periodicamente"

[Timer]
Persistent=True
OnCalendar=*-*-* 02:00:00

[Install]
WantedBy=multi-user.target
